# hostlist.conf
#!upstart

description "Upstart Service Config for hostlist"

start on startup
stop on shutdown
respawn

script
    exec /usr/bin/python2 /srv/hostlist/hostlist.py
end script
